# net/tollgate-wrt/files/etc/uci-defaults/90-tollgate-captive-portal-symlink
#!/bin/sh

# TollGate Captive Portal Symlink Migration
echo "Setting up TollGate NoDogSplash symlink..."

# Check if /etc/nodogsplash/htdocs is already a symlink
if [ -L "/etc/nodogsplash/htdocs" ]; then
    echo "Symlink already exists at /etc/nodogsplash/htdocs"
else
    # If it's a directory, back it up and remove it
    if [ -d "/etc/nodogsplash/htdocs" ]; then
        echo "Backing up existing /etc/nodogsplash/htdocs to /etc/nodogsplash/htdocs.backup"
        mv /etc/nodogsplash/htdocs /etc/nodogsplash/htdocs.backup
    fi
    
    # Remove any existing file/directory at the target location
    rm -rf /etc/nodogsplash/htdocs
    
    # Create the symlink
    echo "Creating symlink from /etc/nodogsplash/htdocs to /etc/tollgate/tollgate-captive-portal-site"
    ln -sf /etc/tollgate/tollgate-captive-portal-site /etc/nodogsplash/htdocs
fi

# No service restart here: uci-defaults runs before procd starts services,
# so nodogsplash will pick up the symlink on its first start.
echo "TollGate NoDogSplash symlink setup completed"

exit 0